# bundle_specs/control-center.app.v1.7.1/payload/service/control-center.service.template
[Unit]
Description=Raspi Homelab Control Center
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={app_dir}
Environment=HOMELAB_ENV_FILE={env_file}
Environment=CONTROL_CENTER_BIND=127.0.0.1
Environment=CONTROL_CENTER_PORT=9000
Environment=HOMELABCTL_BIN={repo_homelabctl}
Environment=REPO_ROOT={repo_root}
ExecStart={python} -m control_center_app
Restart=always
RestartSec=3
User={user}

[Install]
WantedBy=multi-user.target
